# package/etc/service/koruza-alignment/run
#!/bin/bash

cd /koruza/controllers
exec /usr/bin/python /koruza/controllers/alignment.py 2>&1
