# package/etc/service/koruza-vpn/run
#!/bin/bash -e

exec /usr/sbin/openvpn --config /koruza/config/openvpn.conf 2>&1
