# package/etc/service/koruza-nodewatcher/run
#!/bin/bash

exec /koruza/bin/koruza-nodewatcher 2>&1

