# package/etc/service/koruza-reverse-ssh/run
#!/bin/bash

exec /koruza/bin/reverse-ssh 2>&1

