# package/etc/service/koruza-sensors/run
#!/bin/bash

exec /koruza/bin/koruza-sensors 2>&1

