# package/etc/service/koruza-webcam/run
#!/bin/bash

cd /koruza/controllers
exec /usr/bin/python /koruza/controllers/webcam.py 2>&1

