# package/etc/service/koruza-gpio/run
#!/bin/bash

exec /koruza/bin/koruza-gpio 2>&1

