# package/etc/service/koruza-spiral-scan/run
#!/bin/bash

cd /koruza/controllers
exec /usr/bin/python /koruza/controllers/spiral_scan.py 2>&1

