# package/etc/service/koruza-webui/run
#!/bin/bash

cd /koruza/webui
exec /usr/bin/python /koruza/webui/app.py 2>&1

